# root/etc/uci-defaults/luci-zerotier
#!/bin/sh

# ucitrack 只对非 PROCD 程序有效
# NOTE: 注册 firewall include 时使用独立的名字 `firewall.zerotier_include`，
# 不能用 `firewall.zerotier`，否则 init.d/luci-zerotier 在创建 zone 时会
# `uci -q delete firewall.zerotier` 把 include 自己删掉（include 与 zone 同名冲突）。
uci -q batch <<-EOF >/dev/null
	delete ucitrack.@zerotier[-1]
	add ucitrack zerotier
	set ucitrack.@zerotier[-1].init='luci-zerotier'
	commit ucitrack

	delete firewall.zerotier_include
	set firewall.zerotier_include=include
	set firewall.zerotier_include.type=script
	set firewall.zerotier_include.path='/usr/share/zerotier/firewall.include'
	set firewall.zerotier_include.reload=1
	commit firewall
EOF

# Set default persistent config path if not already configured
if ! uci -q get zerotier.global.config_path >/dev/null; then
    uci set zerotier.global.config_path='/etc/zerotier'
    uci set zerotier.global.copy_config_path='1'
    uci commit zerotier
fi

# Seed local.conf from template if neither exists
if [ ! -f /etc/zerotier/local.conf ] && [ ! -f /etc/zerotier-one/local.conf ]; then
    if [ -f /etc/zerotier/local.conf.template ]; then
        # Strip JSON comments (lines starting with //) to create valid JSON
        grep -v '^[[:space:]]*//' /etc/zerotier/local.conf.template | grep -v '^$' > /etc/zerotier/local.conf
    fi
fi

# Set up periodic sync (every hour, saves runtime state to persistent storage)
( crontab -l 2>/dev/null | grep -v 'zerotier-sync'; echo '0 * * * * /usr/bin/zerotier-sync.sh >/dev/null 2>&1' ) | crontab -

# Ensure scripts have execute permissions
# firewall.include MUST be executable: fw3 invokes include scripts via execve()
# when type=script and reload=1; without +x the include is silently skipped.
chmod +x /usr/bin/zerotier-sync.sh 2>/dev/null
chmod +x /usr/share/zerotier/firewall.include 2>/dev/null
[ -f /usr/libexec/rpcd/luci-zerotier ] && chmod +x /usr/libexec/rpcd/luci-zerotier

[ -x /etc/init.d/rpcd ] && /etc/init.d/rpcd restart
[ -x /etc/init.d/uhttpd ] && /etc/init.d/uhttpd restart

rm -f /tmp/luci-indexcache
exit 0
