# root/etc/init.d/luci-zerotier
#!/bin/sh /etc/rc.common

START=99

PROG=/etc/init.d/zerotier

get_config() {
	config_get_bool enabled $1 enabled 0
	config_get_bool nat $1 nat 0
}

log() {
	echo "$@"
	logger -t luci-zerotier "$@" # debug
}

validate_path() {
	case "$1" in
		/etc/zerotier|/etc/zerotier/*|/var/lib/zerotier|/var/lib/zerotier/*|/tmp|/tmp/*)
			return 0
			;;
		*)
			log "Invalid path rejected: $1"
			return 1
			;;
	esac
}

start() {
	config_load zerotier
	config_foreach get_config zerotier

	[ "$enabled" -eq 0 ] && log "Service disabled in config." && return 0

	# Server-side path validation (defense-in-depth, mirrors general.js regex)
	_local_conf=$(uci -q get zerotier.global.local_conf 2>/dev/null)
	_config_path=$(uci -q get zerotier.global.config_path 2>/dev/null)
	[ -n "$_local_conf" ] && { validate_path "$_local_conf" || return 1; }
	[ -n "$_config_path" ] && { validate_path "$_config_path" || return 1; }

	if ! $PROG running 2>/dev/null; then
		log "ZeroTier service is not running!"
		return 1
	fi

	if ! ps 2>/dev/null | grep -q '[/]usr/bin/zerotier-one'; then
		log "ZeroTier daemon (zerotier-one) is not running!"
		return 1
	fi

	[ "$nat" -eq 0 ] && log "NAT disabled in config." && return 0

	local zt_devs zt_timeout zt_count

	# ZeroTier 启动后联网需要时间，每 2 秒检测一次，最多 10 次（共 20 秒）
	zt_timeout=2
	zt_count=0
	zt_max=10

	# 等待真实的 zt 网络设备就绪：通过 JSON 输出提取 portDeviceName
	# 只接受以 "zt" 开头的设备名，过滤掉未就绪时的 null/"-" 等值
	while [ -z "$(zerotier-cli -j listnetworks 2>/dev/null | jsonfilter -e '$[*].portDeviceName' 2>/dev/null | grep '^zt')" ]; do
		[ $zt_count -ge $zt_max ] && {
			log "Timeout waiting for zt device after $((zt_timeout * zt_max))s"
			return 1
		}
		zt_count=$((zt_count + 1))
		sleep $zt_timeout
	done

	zt_devs="$(zerotier-cli -j listnetworks 2>/dev/null | jsonfilter -e '$[*].portDeviceName' 2>/dev/null | grep '^zt')"

	# Configure firewall using uci
	uci -q delete firewall.zerotier
	uci set firewall.zerotier=zone
	uci set firewall.zerotier.name='zerotier'
	uci set firewall.zerotier.input='ACCEPT'
	uci set firewall.zerotier.output='ACCEPT'
	uci set firewall.zerotier.forward='ACCEPT'
	uci set firewall.zerotier.masq='1'
	
	for i in ${zt_devs}; do
		uci add_list firewall.zerotier.device="$i"
	done
	
	# Add forwardings: zerotier <-> lan, zerotier -> wan
	# zerotier -> lan: 让 zerotier 客户端访问局域网
	# lan -> zerotier: 让局域网设备通过 zerotier 互访
	# zerotier -> wan: 让 zerotier 客户端通过本机 WAN 上网（网关模式）
	uci -q delete firewall.zerotier_lan
	uci set firewall.zerotier_lan=forwarding
	uci set firewall.zerotier_lan.src='zerotier'
	uci set firewall.zerotier_lan.dest='lan'

	uci -q delete firewall.lan_zerotier
	uci set firewall.lan_zerotier=forwarding
	uci set firewall.lan_zerotier.src='lan'
	uci set firewall.lan_zerotier.dest='zerotier'

	uci -q delete firewall.zerotier_wan
	uci set firewall.zerotier_wan=forwarding
	uci set firewall.zerotier_wan.src='zerotier'
	uci set firewall.zerotier_wan.dest='wan'
	
	uci commit firewall
	/etc/init.d/firewall reload
	
	log "Added uci firewall rules for ZeroTier interfaces"
}

reload_service() {
	start
}

stop() {
	/usr/bin/zerotier-sync.sh 2>/dev/null

	uci -q delete firewall.zerotier
	uci -q delete firewall.lan_zerotier
	uci -q delete firewall.zerotier_lan
	uci -q delete firewall.zerotier_wan
	uci commit firewall
	/etc/init.d/firewall reload

	log "Synced config and removed uci firewall rules for ZeroTier"
}
